# build/broker/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=mongod --bind_ip 127.0.0.1
user=mongodb
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:apiserver]
command=/usr/bin/cwman api-server
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
