# build/platform/supervisord.conf
[supervisord]
nodaemon=true

[program:bitbucket]
command=%(ENV_BITBUCKET_INSTALL_DIR)s/bin/start-webapp.sh -fg
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:install-plugin]
command=%(ENV_CLOUDWAY_ROOT)s/sbin/install-plugin.sh
stdout_logfile=/dev/null
stderr_logfile=/dev/null
autorestart=false

[program:mongodb]
user=mongodb
command=mongod --bind_ip 127.0.0.1
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:postfix]
command=%(ENV_CLOUDWAY_ROOT)s/sbin/postfix.sh
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
stdout_logfile=/dev/null
stderr_logfile=/dev/null
autorestart=true

[program:hipache]
command=/hipache/bin/start-hipache.sh
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:redis]
user=redis
command=/usr/bin/redis-server /etc/redis/redis.conf
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:apiserver]
command=/usr/bin/cwman api-server
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:proxy-updater]
command=/usr/bin/cwman update-proxy
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:sshd]
command=/usr/bin/cwman sshd
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
